# conf/supervisord_pymicroservice.conf
[supervisord]
nodaemon=true
logfile = /var/log/supervisor.log


[program:uwsgi]
command=/opt/pymicroservice/app/bin/uwsgi --ini  /opt/pymicroservice/conf/uwsgi/uwsgi.ini

